# deploy/systemd/tg-prox-relay.service
[Unit]
Description=TG-prox HTTPS relay server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/TG-prox
EnvironmentFile=/etc/tg-prox-relay.env
ExecStart=/usr/bin/node /opt/TG-prox/server/relay-server.js
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
